# config/crontab
*/1 * * * * /bin/bash /home/ubuntu/jools/scripts/every_minute.sh
5 * * * * /bin/bash /home/ubuntu/jools/scripts/every_hour.sh
15 0 * * * /bin/bash /home/ubuntu/jools/scripts/every_day.sh